# code/Cam.conf
description "Starts Cam Program"
start on startup
task
exec /home/pi/startCam.sh
